# myCgi2.service.in
[Unit]
Description=test

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/@CMAKE_PROJECT_NAME@
#Type=dbus
#BusName=test.service

[Install]
WantedBy=multi-user.target
